# Library/LaunchAgents/uk.tonkin.Keyboard.CalculatorF16.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>uk.tonkin.Keyboard.CalculatorF16</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/hidutil</string>
        <string>property</string>
        <string>--set</string>
        <string>{"UserKeyMapping":[
            {
                "HIDKeyboardModifierMappingSrc": 0xC00000192,
                "HIDKeyboardModifierMappingDst": 0x70000006B
            },
            {
                "HIDKeyboardModifierMappingSrc": 0x700000071,
                "HIDKeyboardModifierMappingDst": 0xC000000CF
            },
            {
                "HIDKeyboardModifierMappingSrc": 0x700000072,
                "HIDKeyboardModifierMappingDst": 0x10000009B
            },
            {
                "HIDKeyboardModifierMappingSrc": 0x700000073,
                "HIDKeyboardModifierMappingDst": 0xC0000029D
            }
        ]}</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
